# framework/crontab
# Cron jobs for master container

* * * * * /bin/sh /repo/framework/cron-1.sh