# scripts/crontab
# Cron jobs for master container

* * * * * /bin/bash /repo/scripts/cron-1.sh