# scripts/crontab
# Cron jobs for master container

* * * * * /bin/sh /repo/scripts/cron-1.sh